# crontab
* 6 * * * python3 /code/scraper.py -c /config.json
* 16 * * * python3 /code/scraper.py -c /config.json
* 22 * * * python3 /code/scraper.py -c /config.json -d 1
